# for-sale-report.timer
[Unit]
Description=Daily timer for For Sale Report

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
