# for-sale-report.service
[Unit]
Description=Run the For Sale Report Application

[Service]
Type=oneshot
ExecStart=/bin/for-sale-report
